# scripts/qa-runner/deploy/qa-runner.service
[Unit]
Description=Vimeflow QA runner daemon
Documentation=https://github.com/winoooops/vimeflow/blob/wip/linear-wiring/scripts/qa-runner/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vimeflow/repo
EnvironmentFile=/etc/vimeflow/qa-runner/control.env
ExecStart=/usr/bin/node /opt/vimeflow/repo/scripts/qa-runner/daemon.js
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=75

[Install]
WantedBy=multi-user.target
